# LFSPkgBuildScripts/net/wpa-supplicant/postinstall
#!/bin/bash

. /LFSPkgBuilds/LFSScripts/SystemData.local

update-desktop-database 2>/dev/null || true

USESYSTEMD=1

if [ $USESYSTEMD -eq 0 ];then
	cat > /etc/sysconfig/ifconfig.wifi0 << EOF
ONBOOT="yes"
IFACE="wlan0"
SERVICE="wpa"

# Additional arguments to wpa_supplicant
WPA_ARGS=""

WPA_SERVICE="ipv4-static"
IP="$IP"
GATEWAY=$GATEWAY
PREFIX="24"
BROADCAST="192.168.1.255"
EOF
else
	cat > /etc/systemd/network/10-wifi-static.network << EOF
[Match]
Name=wlan0
MACAddress=$(cat /sys/class/ieee80211/*/device/net/wlan0/address)

[Network]
Address=${IP}/24
Gateway=${GATEWAY}
DNS=8.8.8.8
EOF
	systemctl enable wpa_supplicant@wlan0||true
	systemctl start wpa_supplicant@wlan0||true
fi
